# LUNA-os/config/includes.chroot/etc/xdg/autostart/blueman.desktop
[Desktop Entry]
Type=Application
Name=Bluetooth Manager
Exec=blueman-applet
NoDisplay=true
X-GNOME-Autostart-enabled=true
