# LUNA-os/config/includes.chroot/etc/xdg/autostart/luna-fix-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Luna Wallpaper Fix
Exec=luna-fix-wallpaper
NoDisplay=true
X-GNOME-Autostart-enabled=true
