# LUNA-os/config/includes.chroot/etc/xdg/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=NetworkManager Applet
Exec=nm-applet
NoDisplay=true
X-GNOME-Autostart-enabled=true
